# configs/config-1.27/systemd/ocne-reset.service
[Unit]
Description=Reset OCNE

[Service]
Type=oneshot
ExecStart=/etc/ocne/ocne-reset.sh
RemainAfterExit=no
