# configs/common/fs/etc/profile.d/harden.sh
# shellcheck shell=sh
# Initialization script for bash, sh, mksh and ksh

umask 0022
TMOUT=900
readonly TMOUT
export TMOUT
