# configs/config-1.26/systemd/update-ca-trust.service
[Unit]
Description="Update CA trust store when new anchors are added"

[Service]
ExecStart=update-ca-trust

[Install]
WantedBy=multi-user.target
